# systemd/user/usb-disconnect-sound.service
[Unit]
Description=Play USB Disconnect Sound
After=sound.target

[Service]
Type=oneshot
# Source environment file if it exists and has valid format
EnvironmentFile=-/etc/xdg/labwc/environment
# Also ensure DISPLAY is set for notifications
Environment="DISPLAY=:0"
ExecStart=/bin/bash -c 'source /etc/xdg/labwc/environment 2>/dev/null || true; /opt/system/lib/play-usb-disconnect-sound.sh'

